# bluebird-camera.service
[Unit]
Description=Bluebird Camera Classifier
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/bluebird/run-bluebird.sh

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=600
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bluebird-camera

[Install]
WantedBy=multi-user.target
